# rds_play.service
[Unit]
Description=rds signal playout

[Service]
WorkingDirectory=/opt/PiFmRds
ExecStart=./rds_play
Type=simple

[Install]
WantedBy=default.target
