# rds_wuvt.service
[Unit]
Description=metadata fetch and rds signal generation
After=network.target

[Service]
WorkingDirectory=/opt/PiFmRds
ExecStart=./rds_wuvt
Type=simple

[Install]
WantedBy=default.target
