# rds_play.service
[Unit]
Description=rds signal playout

[Service]
ExecStart=/opt/PiFmRds/rds_play
Type=simple

[Install]
WantedBy=default.target
